# config/includes.chroot/etc/systemd/system/anvar.service
[Unit]
Description=Anvar Enterprise License Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/anvar/anvar.py
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/anvar.pid
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target

